# supervisord.conf
[supervisord]
nodaemon=true

[program:pootleworker]
command=pootle rqworker

[program:pootle]
command=pootle runfcgi host=0.0.0.0 port=8080 daemonize=false